# ops/container/dropins/hermes-runtime.service.d/llm-timeout.conf
# Patient LLM timeouts so a SLOW or CONTENDED model (e.g. a shared vLLM with
# max-num-seqs=1, where our request queues behind other workloads) is waited on
# instead of aborted at ~60s + retry-stormed. Agentic tasks can be long; we never
# want to kill a request that's merely queued. Nous reads these from the env.
[Service]
Environment=HERMES_STREAM_STALE_TIMEOUT=900
Environment=HERMES_STREAM_READ_TIMEOUT=900
Environment=HERMES_API_TIMEOUT=1800
