# ops/agents-os-edition/gdm/gdm.service.d/10-hermes-agent-gate.conf
# 10-hermes-agent-gate.conf — T022 — drop-in para gdm.service (CTRL-P1-16, G3)
#
# Instalar en: /etc/systemd/system/gdm.service.d/10-hermes-agent-gate.conf
#
# PROPÓSITO: impedir que GDM arranque (y por tanto que el usuario entre al
# escritorio) antes de que el agente haya notificado READY=1.
#
# La dependencia pasa por agents-os-healthy.target (no directamente por
# hermes-runtime-ready.target) porque agents-os-healthy.target lleva el
# JobTimeoutSec=300 + JobTimeoutAction=reboot que produce el OTA-rollback
# si el sistema no arranca en 5 minutos. GDM simplemente espera a ese target.
#
# Relación entre targets:
#   hermes-runtime.service (READY=1)
#       → hermes-runtime-ready.target
#           → agents-os-healthy.target   ← ESTE drop-in lo consume
#               → gdm.service   ← ESTE drop-in lo ancla aquí
#
# Wants (no Requires): si agents-os-healthy.target no existe en la imagen
# (boot en hardware sin el perfil agents-os), GDM arranca igual —
# degradación segura para entornos de desarrollo sin los targets instalados.
# En producción agents-os-healthy.target SIEMPRE existe.
#
# After: asegura el orden de inicio (Wants solo garantiza activación,
# no orden). Sin After, systemd podría arrancar GDM en paralelo con el target.

[Unit]
After=agents-os-healthy.target
Wants=agents-os-healthy.target
