# ops/container/dropins/hermes-browser-netns.service.d/container.conf
# Container delta: BOTH `ip netns exec` AND `ip -n NETNS` remount /sys on netns
# entry — which the tailored seccomp profile (safent.json) DENIES (no @mount) →
# "mount of /sys failed: Operation not permitted" → the netns-side veth never gets
# its IP/route → the browser/terminal cannot reach the egress proxy (10.200.0.1)
# → no domain reachable. Fix: `nsenter --net=<path>` enters ONLY the network
# namespace (no mount-ns change, no /sys remount) — `ip` uses netlink, not /sys, so
# addr/link/route apply cleanly under the mount-denying seccomp. Same jail, exit 0.
# (The OS/VM build keeps `ip netns exec`; this overrides ONLY in-container.)
[Service]
ExecStart=
ExecStart=/bin/sh -c 'ip netns add hermes-browser 2>/dev/null || true'
ExecStart=/bin/sh -c 'ip link add veth-hbr-host type veth peer name veth-hbr-ns 2>/dev/null || true && ip link set veth-hbr-ns netns hermes-browser 2>/dev/null || true && ip addr add 10.200.0.1/30 dev veth-hbr-host 2>/dev/null || true && ip link set veth-hbr-host up'
ExecStart=/bin/sh -c 'N="nsenter --net=/run/netns/hermes-browser"; $N ip addr add 10.200.0.2/30 dev veth-hbr-ns 2>/dev/null || true; $N ip link set veth-hbr-ns up; $N ip link set lo up; $N ip route add default via 10.200.0.1 2>/dev/null || true'
ExecStart=/bin/sh -c 'sysctl -w net.ipv4.conf.veth-hbr-host.forwarding=1 2>/dev/null || true'
ExecStart=/bin/sh -c 'nft -f /etc/nftables/hermes/browser-host.nft'
ExecStart=/bin/sh -c 'nsenter --net=/run/netns/hermes-browser nft -f /etc/nftables/hermes/browser-ns.nft'

# ── hermes-mcp netns (container delta — nsenter form) ─────────────────────────
# The main unit (hermes-browser-netns.service:138-195) builds the hermes-mcp netns
# with `ip netns exec` which remounts /sys on entry → DENIED by the container seccomp
# (no @mount) → the MCP netns/veth never come up → the egress proxy's second listener
# (10.200.1.1:3128) can't bind → MCP dead. The `ExecStart=` reset above silently
# dropped the whole hermes-mcp block; re-add it here in nsenter form (network-ns only,
# no /sys remount), mirroring the browser steps. Same jail, exit 0.
# (Regression caught by the 2026-06-19 integrated re-validation: secure-but-broken.)
ExecStart=/bin/sh -c 'ip netns add hermes-mcp 2>/dev/null || true'
ExecStart=/bin/sh -c 'ip link add veth-hmcp-host type veth peer name veth-hmcp-ns 2>/dev/null || true && ip link set veth-hmcp-ns netns hermes-mcp 2>/dev/null || true && ip addr add 10.200.1.1/30 dev veth-hmcp-host 2>/dev/null || true && ip link set veth-hmcp-host up'
ExecStart=/bin/sh -c 'N="nsenter --net=/run/netns/hermes-mcp"; $N ip addr add 10.200.1.2/30 dev veth-hmcp-ns 2>/dev/null || true; $N ip link set veth-hmcp-ns up; $N ip link set lo up; $N ip route add default via 10.200.1.1 2>/dev/null || true'
ExecStart=/bin/sh -c 'sysctl -w net.ipv4.conf.veth-hmcp-host.forwarding=1 2>/dev/null || true'
ExecStart=/bin/sh -c 'nft -f /etc/nftables/hermes/mcp-host.nft'
ExecStart=/bin/sh -c 'nsenter --net=/run/netns/hermes-mcp nft -f /etc/nftables/hermes/mcp-ns.nft'
ExecStart=/bin/sh -c 'if nft list table inet hermes_host >/dev/null 2>&1; then if ! nft list chain inet hermes_host input 2>/dev/null | grep -q "mcp netns -> egress proxy"; then nft add rule inet hermes_host input ip saddr 10.200.1.2 ip daddr 10.200.1.1 tcp dport 3128 accept comment "mcp netns -> egress proxy"; fi; fi'
