# ops/container/dropins/hermes-shell-server.service.d/container.conf
# Container delta: bind the UI on all interfaces INSIDE the container so the
# `-p 127.0.0.1:17517:7517` publish reaches it. The host port mapping is the new
# edge (replaces the old 127.0.0.1 + VZ relay). The OS unit defaulted to
# 127.0.0.1. NOTE: 0.0.0.0 is REQUIRED for podman/docker port-publish — the host
# only sees 127.0.0.1:PORT, the container side must listen on all interfaces or
# the publish has nothing to forward to. Binding 127.0.0.1 here breaks the UI.
# The HTTP edge is protected by token auth, not by the bind address.
#
# C3 PASS-5 (bootstrap secret uid-decoupling + UNMOVABLE, CLASS fix). PASS-1 wrote
# the secret 0600 owned by uid 880 (hermes). PASS-2 split it into root:root
# plaintext + commitment. PASS-3 moved both into a bootstrap/ subdir (root:root).
# PASS-4 recreated that subdir from scratch every boot. But through PASS-4 the
# PARENT was STILL /var/lib/hermes, owned 0755 by uid 880 (hermes — the daemon,
# the named adversary). On Unix the right to unlink/rename a directory ENTRY is
# governed by WRITE on the CONTAINING dir, not the entry's mode/owner — so AFTER
# provisioning, uid 880 could `rename` the whole root-owned `bootstrap` dir entry
# away and drop in its OWN dir/symlink holding a commitment of a secret it chose;
# the reader opened by plain path and would follow + trust it (red-team PASS-5
# 2026-06-19).
#
# CLASS FIX — two independent legs (either alone shuts the door):
#
#   1. The bootstrap dir lives under a SEPARATE, ROOT-OWNED parent that uid
#      880/886 cannot write — /var/lib/hermes-bootstrap (0755 root:root), a
#      SIBLING of /var/lib/hermes, created by tmpfiles (ops/.../tmpfiles/
#      hermes.conf). uid 880 cannot rename/replace the `bootstrap` entry because
#      it cannot write that parent at all.
#
#   2. The reader (_load_bootstrap_commitment) O_NOFOLLOW-opens the commitment and
#      fstat-verifies the FILE and its DIR are root-owned before trusting the
#      bytes — a swapped entry (symlink/non-root file/non-root dir) is rejected,
#      fail-closed.
#
#   /var/lib/hermes-bootstrap/                          0755 root:root  (tmpfiles)
#   /var/lib/hermes-bootstrap/bootstrap/                0755 root:root  (per boot)
#       Recreated fresh each boot. Only root can create/rename/unlink entries.
#
#   …/bootstrap/webui-bootstrap          PLAINTEXT  0400 root:root
#       Readable ONLY by root. The legitimate owner reads it from the HOST
#       authenticated session — `podman exec` runs as root in the container:
#           podman exec safent-os cat /var/lib/hermes-bootstrap/bootstrap/webui-bootstrap
#       then opens the UI ONCE at  http://localhost:17517/?k=<that-secret>
#
#   …/bootstrap/webui-bootstrap.commit   COMMITMENT 0444 root:root
#       Hex SHA-256 of the secret. World-readable so the shell-server (uid 880)
#       can VERIFY a presented secret, root-only-writable, non-invertible, and —
#       critically — UNMOVABLE: a compromised uid-880/uid-886 process can neither
#       overwrite it (not owner) nor rename/replace it (cannot write the
#       root-owned parent /var/lib/hermes-bootstrap nor the bootstrap subdir).
#
# Provisioned by a ROOT ExecStartPre (the `+` prefix runs as root despite
# User=hermes; works under NoNewPrivileges= because PID1 spawns it directly as
# uid 0) BEFORE the shell-server (User=hermes) starts. The shell-server only
# READS the commitment, never writes either file. `_provision_bootstrap_secret()`
# O_NOFOLLOW-opens the root parent, fstat-verifies it is root-owned, rm -rf any
# existing bootstrap entry, mkdir fresh root:root 0755, fstat-verify it, then
# writes a FRESH secret (0400) + commitment (0444) via openat O_EXCL|O_NOFOLLOW.
# No early-exit, no symlink following, no TOCTOU window.
#
# Tradeoff: the secret rotates on every boot. The owner re-reads it after each
# boot (never a bookmarked URL) — an unforgeable gate beats a poisonable one:
#   podman exec safent-os cat /var/lib/hermes-bootstrap/bootstrap/webui-bootstrap
# On success GET / mints a SHORT-LIVED ROTATING session token (not the master
# mutator key). Override paths / TTL via HERMES_SHELL_BOOTSTRAP /
# HERMES_SHELL_BOOTSTRAP_COMMIT / HERMES_SHELL_SESSION_TTL_S if needed.
[Service]
Environment=HERMES_SHELL_BIND_HOST=0.0.0.0
Environment=HERMES_SHELL_BOOTSTRAP=/var/lib/hermes-bootstrap/bootstrap/webui-bootstrap
Environment=HERMES_SHELL_BOOTSTRAP_COMMIT=/var/lib/hermes-bootstrap/bootstrap/webui-bootstrap.commit
Environment=HERMES_SHELL_SESSION_TTL_S=3600
ExecStartPre=+/usr/bin/python3 -c "from hermes.shell_server.main import _provision_bootstrap_secret; _provision_bootstrap_secret()"
