# ops/agents-os-edition/systemd/bootc-updater.timer
# bootc-updater.timer — comprueba canal OTA cada 12h con jitter.

[Unit]
Description=Hermes OTA Updater periodic check

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=4h
Persistent=yes

[Install]
WantedBy=timers.target
