# ops/agents-os-edition/systemd/hermes-config-sync.service
# hermes-config-sync.service — Fase 4: cloud→associate policy sync.
#
# Pulls a signed PolicyBundle from the cloud endpoint, verifies the
# Ed25519 signature, and applies the policy via the local D-Bus verbs.
#
# Enabled UNCONDITIONALLY in every image (the Containerfile systemctl-enables it
# for both editions). It stays inert until the instance is paired: the
# ConditionPathExists gate below blocks start until the .enterprise marker exists
# (written by `lumen pair`), and config_sync itself exits early if
# not is_associated(). So a plain community image becomes the associate the
# moment it is paired — no edition-specific build step is involved.
#
# To enable manually for testing inside the container:
#   systemctl enable --now hermes-config-sync.service

[Unit]
Description=Hermes config sync — cloud policy pull (Enterprise Edition)
# The runtime daemon MUST be ready before we push D-Bus mutators.
After=network-online.target hermes-runtime.service hermes-keygen.service
Wants=network-online.target hermes-keygen.service
Requires=hermes-runtime.service
PartOf=hermes-workspace.target
# Activation gate WITHOUT a bespoke ops script: the unit is enabled by default
# (see [Install]) but only RUNS when the instance is paired. The enterprise
# marker is written to the persistent volume on `lumen pair` and removed on
# `lumen unpair`, so this condition makes the sync a no-op in community edition
# (no restart-loop) and auto-starts it on every boot once paired. The loop also
# self-exits if is_associated() is false (defence in depth).
ConditionPathExists=/var/lib/hermes/instance/.enterprise

[Service]
Type=simple
User=hermes
Group=hermes
# HERMES_HOME must be writable for hermes_cli config reads (same as other units).
Environment=HERMES_HOME=/var/lib/hermes/hermes-home
Environment=COMPOSIO_CACHE_DIR=/var/lib/hermes/composio-cache
# The state DB path can be overridden for testing; default matches shell-server.
Environment=HERMES_STATE_DB=/var/lib/hermes/shell-state.db
# Sync interval in seconds (default 300).  Override with a drop-in for fast testing.
Environment=HERMES_CONFIG_SYNC_INTERVAL=300

ExecStart=/usr/bin/python3 -m hermes.config_sync
Restart=always
RestartSec=10

Slice=agents-os.slice
MemoryMax=64M
OOMScoreAdjust=200

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
RestrictNamespaces=yes
LockPersonality=yes
ReadWritePaths=/var/lib/hermes
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
# Enabled by default in the image (Containerfile: systemctl enable). The
# ConditionPathExists gate above keeps it inert until the instance is paired,
# so there is no separate "Fase 6 activation script" to forget — pairing writes
# the marker and the next boot (or `lumen pair` starting it) brings sync up.
WantedBy=hermes-workspace.target
