# ops/agents-os-edition/systemd/hermes-remote-vnc-password-gen.service
# hermes-remote-vnc-password-gen.service — genera /etc/hermes/remote-vnc-password
# en el PRIMER BOOT.
#
# SECURITY FIX (finding #2 — hardcoded fallback "hermesrc"):
# hermes-grd-init usaba `|| echo "hermesrc"` como fallback si el archivo no existía.
# Eso significa que CUALQUIER instancia sin este archivo usaba una contraseña conocida
# y pública. Esta unit genera el archivo antes de que hermes-grd-init.service arranque,
# de modo que el fallback nunca se activa; hermes-grd-init verifica el archivo o falla.
#
# Idempotente: ConditionPathExists=! evita regenerar en reboots (regenerar la pass
# invalida el llavero GRD que la almacenó en el arranque previo).
# Must run BEFORE hermes-grd-init.service (user session).
#
# VNC password constraint: grdctl acepta máximo 8 caracteres (límite del protocolo
# RFB). Se generan exactamente 8 bytes del espacio base-58 para máxima entropía
# dentro del límite del protocolo.

[Unit]
Description=Hermes — Remote VNC password generation (first boot)
After=systemd-tmpfiles-setup.service local-fs.target
# La unit de usuario hermes-grd-init.service arranca dentro de la sesión GDM
# (graphical-session.target). Asegurar que esta unit esté Before=graphical.target
# garantiza que el archivo exista antes de que GDM autologuee a hermes-user.
Before=graphical.target display-manager.service
ConditionPathExists=!/etc/hermes/remote-vnc-password
RequiresMountsFor=/etc/hermes

[Service]
Type=oneshot
RemainAfterExit=yes
User=root
Group=root

# Genera exactamente 8 caracteres alfanuméricos sin ambiguos desde /dev/urandom.
# 8 chars = límite máximo del protocolo VNC/RFB (RFC 6143 §7.2.2).
# 0600 root:root — solo root puede leerlo; hermes-grd-init lo lee como root via
# el servicio de usuario que arranca con hermes-user pero necesita acceso a /etc.
# NOTA: hermes-grd-init.service es un servicio de USUARIO (no root), por eso
# el ownership es root:hermes y permisos 0640 para que el proceso hermes-user lo lea.
ExecStart=/bin/sh -c '\
    mkdir -p /etc/hermes && \
    tr -dc "A-HJ-NP-Za-km-z2-9" < /dev/urandom | head -c 8 > /etc/hermes/remote-vnc-password && \
    chown root:hermes-user /etc/hermes/remote-vnc-password && \
    chmod 0640 /etc/hermes/remote-vnc-password'

NoNewPrivileges=yes
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
LockPersonality=yes
ReadWritePaths=/etc/hermes

[Install]
WantedBy=multi-user.target
