# ops/container/dropins/hermes-exec-launcher.service.d/macos-cgroup.conf
[Service]
# podman-machine (macOS) compat: the Mac podman VM cannot remount /sys/fs/cgroup read-only
# inside a unit's mount namespace (ProtectControlGroups=yes) → "Failed to set up mount
# namespacing: /sys/fs/cgroup: Permission denied" → 226/NAMESPACE → the unit crash-loops.
# This is defense-in-depth ONLY; the real Lumen cage (netns/nftables, uid-split, Landlock,
# seccomp) does not depend on it. Disable so the namespace setup succeeds. On native Linux
# this is a small DiD reduction; functionality is unchanged (the daemon manages cgroups via
# systemd APIs, not by writing /sys/fs/cgroup directly).
ProtectControlGroups=no
