# ops/agent-cage/systemd/hermes-runtime.service.d/openshell.conf
# OpenShell confinement (P2) — OPT-IN. The deploy writes /var/lib/hermes/openshell/openshell.env
# with HERMES_OPENSHELL_SANDBOX=<per-agent sandbox> once the gateway + sandbox are provisioned.
# Absent → the terminal falls back to the exec-launcher cage (fail-closed). Never in-daemon.
[Service]
EnvironmentFile=-/var/lib/hermes/openshell/openshell.env
Environment=HERMES_OPENSHELL_CLI=/usr/libexec/hermes/openshell
Environment=HERMES_OPENSHELL_HOME=/var/lib/hermes/openshell
Environment=LD_LIBRARY_PATH=/usr/lib/hermes/openshell
